# etc/udev/rules.d/99-governor.rules
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="~/.local/bin/set-governor.sh"
